# etc/service/lucky/run
#!/bin/bash -e

exec chpst -u lucky:lucky /usr/local/bin/node /lucky/src/js/server.js 2>&1
